# SmartHomePCControl/Linux/shutdown-listener.service
[Unit]
Description=TCP Shutdown Listener
After=network.target

[Service]
ExecStart=/usr/local/bin/shutdown-listener
Restart=always

[Install]
WantedBy=multi-user.target
